# etc/supervisor/conf.d/service.conf
[program:service]
command=/usvc/service-wrapper.sh /usvc/service
directory=/usvc
redirect_stderr=true